# out/artifacts/VPNnetworkCheck_jar/VPNnetworkCheck.service
[Unit]
Description=VPNnetworkCheck
After=syslog.target network.target

[Service]
Type=simple

ExecStart=java -jar /home/vladislav/network-monitor/VPNnetworkCheck.jar

[Install]
WantedBy=multi-user.target
